# tdb-api/deploy/supervisord.conf


[unix_http_server]
file=/tmp/%(ENV_project)s/%(ENV_MODE)s/supervisor.sock   ; (the path to the socket file)


[supervisord]
logfile=/tmp/%(ENV_project)s/%(ENV_MODE)s/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/%(ENV_project)s/%(ENV_MODE)s/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/%(ENV_project)s/%(ENV_MODE)s/supervisor.sock ; use a unix:// URL  for a unix socket


[program:api]
command=%(ENV_INSTALL_DIR)s/venv/bin/gunicorn run:app -b unix:///tmp/%(ENV_project)s/%(ENV_MODE)s/api.sock --threads 80 -w %(ENV_processes)s
directory=%(ENV_INSTALL_DIR)s
environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", MODE="%(ENV_MODE)s"
autostart=true
autorestart=true
stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/api_err_log
stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/api_out_log

[program:commission_engine]
command=%(ENV_INSTALL_DIR)s/venv/bin/python commission_engine.py
directory=%(ENV_INSTALL_DIR)s/app/engine
environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
autostart=true
autorestart=true
stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/commission_engine_err_log
stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/commission_engine_out_log

[program:crawler_engine]
command=%(ENV_INSTALL_DIR)s/venv/bin/python crawler_engine.py
directory=%(ENV_INSTALL_DIR)s/app/engine
environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
autostart=true
autorestart=true
stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/crawler_engine_err_log
stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/crawler_engine_out_log

[program:check_order_engine]
command=%(ENV_INSTALL_DIR)s/venv/bin/python check_order_engine.py
directory=%(ENV_INSTALL_DIR)s/app/engine
environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
autostart=true
autorestart=true
stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/check_order_engine_log
stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/check_order_engine_log

[program:generate_order_engine]
command=%(ENV_INSTALL_DIR)s/venv/bin/python generate_order_engine.py
directory=%(ENV_INSTALL_DIR)s/app/engine
environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
autostart=true
autorestart=true
stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/generate_order_engine_log
stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/generate_order_engine_log

[program:match_order_engine]
command=%(ENV_INSTALL_DIR)s/venv/bin/python match_order_engine.py
directory=%(ENV_INSTALL_DIR)s/app/engine
environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
autostart=true
autorestart=true
stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/match_order_engine_log
stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/match_order_engine_log

[program:confirm_order_engine]
command=%(ENV_INSTALL_DIR)s/venv/bin/python confirm_order_engine.py
directory=%(ENV_INSTALL_DIR)s/app/engine
environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
autostart=true
autorestart=true
stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/confirm_order_engine_log
stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/confirm_order_engine_log

;[program:eth_engine]
;command=%(ENV_INSTALL_DIR)s/venv/bin/python eth_engine.py
;directory=%(ENV_INSTALL_DIR)s/app/engine
;environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
;autostart=true
;autorestart=true
;stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/eth_engine_err_log
;stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/eth_engine_out_log


;[program:dash_engine]
;command=%(ENV_INSTALL_DIR)s/venv/bin/python dash_engine.py
;directory=%(ENV_INSTALL_DIR)s/app/engine
;environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
;autostart=true
;autorestart=true
;stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/dash_engine_err_log
;stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/dash_engine_out_log


;[program:ltc_engine]
;command=%(ENV_INSTALL_DIR)s/venv/bin/python ltc_engine.py
;directory=%(ENV_INSTALL_DIR)s/app/engine
;environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
;autostart=true
;autorestart=true
;stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/ltc_engine_err_log
;stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/ltc_engine_out_log

;[program:usdt_engine]
;command=%(ENV_INSTALL_DIR)s/venv/bin/python usdt_engine.py
;directory=%(ENV_INSTALL_DIR)s/app/engine
;environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
;autostart=true
;autorestart=true
;stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/usdt_engine_err_log
;stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/usdt_engine_out_log

;[program:btc_engine]
;command=%(ENV_INSTALL_DIR)s/venv/bin/python btc_engine.py
;directory=%(ENV_INSTALL_DIR)s/app/engine
;environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
;autostart=true
;autorestart=true
;stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/btc_engine_err_log
;stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/btc_engine_out_log

;[program:wc_engine]
;command=%(ENV_INSTALL_DIR)s/venv/bin/python wc_engine.py
;directory=%(ENV_INSTALL_DIR)s/app/engine
;environment = PATH="%(ENV_INSTALL_DIR)s/venv/bin", PYTHONPATH="%(ENV_INSTALL_DIR)s", MODE="%(ENV_MODE)s"
;autostart=true
;autorestart=true
;stderr_logfile=%(ENV_INSTALL_DIR)s/../logs/wc_engine_err_log
;stdout_logfile=%(ENV_INSTALL_DIR)s/../logs/wc_engine_out_log
